# release/src-rt-5.04behnd.4916/router-sysdep.gt-be19000ai/wlan/rdk/meta-wifi/recipes-ccsp/ccsp/ccsp-common-library/checkbrcmwifisupport.service
[Unit]
Description=Check WIFI support on Broadcom Device

After=broadcom-wifi-drivers.service
Requires=broadcom-wifi-drivers.service

[Service]
Type=forking
ExecStart=/bin/sh -c '/usr/ccsp/wifi/brcm_checkwifi.sh'

StandardOutput=syslog+console
