# release/src/router/wireguard-tools/src/systemd/wg-quick.target
[Unit]
Description=WireGuard Tunnels via wg-quick(8)
