# release/src/router/e2fsprogs/scrub/e2scrub_all.timer.in
[Unit]
Description=Periodic ext4 Online Metadata Check for All Filesystems

[Timer]
# Run on Sunday at 3:10am, to avoid running afoul of DST changes
OnCalendar=Sun *-*-* 03:10:00
RandomizedDelaySec=60
Persistent=true

[Install]
WantedBy=timers.target
