# release/src-rt-5.04behnd.4916/router-sysdep.rt-be86u/wlan/rdk/meta-wifi/recipes-ccsp/ccsp/ccsp-common-library/checkbrcmwifisupport.path
[Unit]
Description=Check Broadcom  Wifi Support Initialized Path

[Path]
PathExists=/tmp/pam_initialized
Unit=checkbrcmwifisupport.service

[Install]
WantedBy=multi-user.target
