# release/src/router/nfs-utils-1.3.4/systemd/proc-fs-nfsd.mount
[Unit]
Description=NFSD configuration filesystem

[Mount]
What=nfsd
Where=/proc/fs/nfsd
Type=nfsd
