# release/src-rt-5.04axhnd.675x/bcmdrivers/broadcom/net/wl/impl87/main/components/rdk/meta-rdk-broadcom-wifi/meta-brcm-wifi/recipes-bsp/broadcom-wifi/files/wifi.service
[Unit]
Description=WiFi Service.

After=brcm-net.service
Requires=brcm-net.service

[Service]
Type=oneshot
ExecStart=/bin/sh -c '/lib/rdk/wifi_setup.sh init'
ExecReload=/bin/sh -c '/lib/rdk/wifi_setup.sh restart'
ExecStop=/bin/sh -c '/lib/rdk/wifi_setup.sh stop'
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
