# release/src/router/strongswan-5.2.1/init/systemd/strongswan.service.in
[Unit]
Description=strongSwan IPsec IKEv1/IKEv2 daemon using ipsec.conf
After=syslog.target

[Service]
ExecStart=@SBINDIR@/@IPSEC_SCRIPT@ start --nofork
StandardOutput=syslog

[Install]
WantedBy=multi-user.target
