# release/src/router/net-snmp-5.7.2/dist/snmpd.service
#
# SNMP agent service file for systemd
#
#
# The service should be enabled, i.e. snmpd should start during machine boot.
# Socket activation shall not be used. See README.systemd for details.

[Unit]
Description=Simple Network Management Protocol (SNMP) daemon.
After=syslog.target network.target

[Service]
# Type=notify is also supported. It should be set when snmpd.socket is not used.
Type=simple
ExecStart=/usr/sbin/snmpd -f

[Install]
WantedBy=multi-user.target
