# release/src/router/sysstat-10.0.3/cron/sysstat.cron.hourly.in
#!/bin/sh
# Run system activity accounting tool every @CRON_INTERVAL@ minutes
@SA_LIB_DIR@/sa1 @CRON_INTERVAL_SEC@ @CRON_COUNT@ &
