# release/src/router/freeradius-server-3.0.0/scripts/radiusd.cron.monthly.in
#! /bin/sh
#
# radiusd	Cron script to rotate radwtmp file monthly.
#

prefix=@prefix@
localstatedir=@localstatedir@
logdir=@logdir@

umask 022
cd $logdir

# Take care of the standard logfiles.
cd $logdir
if [ -f radwtmp ]
then
	savelog -g adm -m 644 -c 6 radwtmp > /dev/null
fi

