# release/src-rt-5.04behnd.4916/bcmdrivers/broadcom/net/wl/impl103/main/components/rdk/meta-rdk-broadcom-wifi/meta-brcm-wifi/recipes-ccsp/ccsp/ccsp-common-library/brcmwifiinitialized.path
[Unit]
Description=BRCM WiFi Ready Path

[Path]
PathExists=/tmp/.brcm_wifi_ready
Unit=ccspwifiagent.service

[Install]
WantedBy=multi-user.target
