# release/src-rt-5.04axhnd.675x/bcmdrivers/broadcom/net/wl/impl87/main/components/rdk/meta-rdk-broadcom-wifi/meta-brcm-wifi/recipes-ccsp/ccsp/ccsp-common-library/brcmwifiinitialized.path
[Unit]
Description=BRCM WiFi Initialized Path
[Path]
PathExists=/tmp/wl_instances.txt
Unit=ccspwifiagent.service
[Install]
WantedBy=multi-user.target
