# release/src-rt-5.04behnd.4916/router-sysdep.gt-be98_pro/wlan/rdk/meta-wifi/recipes-ccsp/ccsp/ccsp-common-library/brcmwifiinitialized.path
[Unit]
Description=BRCM WIFI Initialized Path

[Path]
PathExists=/tmp/.brcm_wifi_ready
Unit=ccspwifiagent.service

[Install]
WantedBy=multi-user.target
