# release/src/router/bluez-5.41/obexd/src/obex.service.in
[Unit]
Description=Bluetooth OBEX service

[Service]
Type=dbus
BusName=org.bluez.obex
ExecStart=@libexecdir@/obexd

[Install]
Alias=dbus-org.bluez.obex.service
