# release/src/router/bluez-5.56/tools/bluetooth-logger.service.in
[Unit]
Description=Bluetooth monitor logger
ConditionPathIsDirectory=/sys/class/bluetooth

[Service]
Type=simple
ExecStart=@pkglibexecdir@/btmon-logger -p -b /var/log/bluetooth/hci.log
NotifyAccess=main
CapabilityBoundingSet=CAP_NET_RAW
LimitNPROC=1
ProtectHome=true
ProtectSystem=full
PrivateTmp=true
PrivateDevices=true

[Install]
WantedBy=bluetooth.target
